# infra/systemd/ah-scraper.service
[Unit]
Description=Atlantic Highlands scheduled document scraper
After=network-online.target ah-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/atlantic-highlands/api
# Inherit ah-api's environment (DATABASE_URL, AWS creds, AHNJ_PROXY_*).
# AHNJ_PROXY_SECRET in particular must be present or ahnj.com requests
# return 403 from the Cloudflare Worker. Drop-in this service must have
# the same proxy override that ah-api uses:
#   sudo mkdir -p /etc/systemd/system/ah-scraper.service.d/
#   sudo cp /etc/systemd/system/ah-api.service.d/ahnj-proxy.conf \
#           /etc/systemd/system/ah-scraper.service.d/ahnj-proxy.conf
#   sudo systemctl daemon-reload
# The drop-in file is not in version control (it carries the secret).
EnvironmentFile=-/etc/atlantic-highlands.env

ExecStart=/opt/atlantic-highlands/api/venv/bin/python -m scripts.scheduled_scrape

# Don't crash the timer if a scrape errors out. Errors get logged via
# journal + persisted to scraper_runs.errors.
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
